# scripts/ppsa-firstboot.service
[Unit]
Description=PPSA First Boot Progress Display
Documentation=man:ppsa
# Start as soon as the system is up and the install service is enabled.
# We do NOT After=ppsa-install.service because we want to start showing
# progress even before the install service has begun its first step.
After=systemd-logind.service getty-pre.target
Before=getty@tty1.service
# Only run on first boot. After install, the .installed flag exists and
# we never run again.
ConditionPathExists=!/opt/ppsa/.installed
# No Conflicts=getty@tty1.service here: systemd honors a unit's Conflicts=
# as part of the boot transaction even when that unit's own ConditionPathExists
# then skips it, which permanently blocked getty@tty1.service from ever
# starting again on every boot after the first (no console, no login prompt,
# looked like a hang). getty@tty1.service instead carries its own
# ConditionPathExists=/opt/ppsa/.installed override (see the
# getty@tty1.service.d/autologin.conf drop-in in build-live-usb.sh), so it
# simply doesn't start at all until install.sh has finished — no fighting
# over the tty needed. TTYReset=yes/TTYVHangup=yes below are kept as a
# belt-and-suspenders reclaim in case anything else is holding tty1.

[Service]
Type=simple
# This service owns /dev/tty1 during first boot. It writes the progress
# UI to the tty and exits when install completes (or the user presses
# a key), at which point getty takes over (autologin as ppsa).
ExecStart=/opt/ppsa/scripts/ppsa-firstboot.sh
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
# Run as root so we can read /var/log and the install flag.
User=root
# Don't restart on failure — if the script crashes, let the user see the
# login prompt and diagnose.
Restart=no
# Give the user 30s to press a key, then exit (handled inside the script).
TimeoutStopSec=35

[Install]
WantedBy=multi-user.target
