# scripts/ppsa-firewall-restore.service
[Unit]
Description=PPSA - Reapply WG_FRIENDS iptables chain at boot
Documentation=https://github.com/kaiser62/ppsa
# Reapply AFTER docker and the compose stack are up. This is critical:
# Docker generates its own nat/FORWARD rules (including the published-port
# DNAT for the WebUI container) fresh on every boot, and those rules point
# at whatever bridge IP the container gets *this* boot. We must not clobber
# them. Older versions did a full `iptables-restore` of a saved snapshot,
# which froze Docker's DNAT at the container's *first-boot* IP; after a
# reboot the container gets a new IP, so the restored DNAT pointed at a dead
# address and the WebUI became unreachable over WireGuard/LAN (only
# localhost + direct container IP still worked). Instead we simply re-run
# the idempotent apply script, which rebuilds ONLY the WG_FRIENDS chain and
# its INPUT jump from firewall.json and leaves Docker's rules untouched.
After=network-online.target docker.service ppsa-docker-compose.service
Wants=network-online.target ppsa-docker-compose.service
# No-op on installs where the apply script isn't present yet.
ConditionPathExists=/opt/ppsa/scripts/ppsa-firewall-apply.sh

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/ppsa/scripts/ppsa-firewall-apply.sh
User=root

[Install]
WantedBy=multi-user.target
